# infra/dogcare-api.service
[Unit]
Description=DogCare API (Node.js)
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=azureuser
WorkingDirectory=/srv/dogcare-api
EnvironmentFile=/srv/dogcare-api/.env
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=3
StandardOutput=append:/var/log/dogcare-api.log
StandardError=append:/var/log/dogcare-api.log

[Install]
WantedBy=multi-user.target
